# qresos-backend.service
[Unit]
Description=QRes OS 4 Restaurant Management System Backend
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/qresos4/backend
Environment=PATH=/var/www/qresos4/backend/venv/bin
Environment=VIRTUAL_ENV=/var/www/qresos4/backend/venv
Environment=NODE_ENV=production
Environment=DEBUG=false
Environment=RELOAD=false
Environment=HOST=192.168.4.1
Environment=PORT=8000
Environment=LOG_LEVEL=info
ExecStartPre=/bin/bash -c 'cd /var/www/qresos4/backend && source venv/bin/activate && python3 -m alembic upgrade head'
ExecStart=/bin/bash /var/www/qresos4/backend/start.sh
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
