# qresos-backend.service
[Unit]
Description=QRes OS 4 Restaurant Management System Backend
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=admin
Group=admin
WorkingDirectory=/home/admin/qresos/backend
Environment=PATH=/home/admin/qresos/backend/venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=VIRTUAL_ENV=/home/admin/qresos/backend/venv
Environment=NODE_ENV=production
Environment=DEBUG=false
Environment=RELOAD=false
Environment=PORT=8000
Environment=LOG_LEVEL=info
ExecStartPre=/bin/bash -c 'cd /home/admin/qresos/backend && if [ -d "venv" ]; then source venv/bin/activate; fi && python3 -m alembic upgrade head'
ExecStart=/bin/bash /home/admin/qresos/backend/start.sh
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
